# vendor/bundle/ruby/2.7.0/gems/unicorn-5.7.0/examples/unicorn.socket
# ==> /etc/systemd/system/unicorn.socket <==
[Unit]
Description = unicorn sockets

[Socket]
ListenStream = 127.0.0.1:8080
ListenStream = /tmp/path/to/.unicorn.sock
Service = unicorn@1.service

[Install]
WantedBy = sockets.target
