# vendor/pkgr/scaling/systemd/openproject-backup-PROCESS_NUM.service
[Unit]
StopWhenUnneeded=true
Requires=openproject-backup.service
After=openproject-backup.service

[Service]
Environment=APP_PROCESS_INDEX=__PROCESS_INDEX__
ExecStart=/usr/bin/openproject run backup
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=%n

[Install]
WantedBy=openproject-backup.service
