# vendor/pkgr/scaling/systemd/openproject-check.service
[Unit]
StopWhenUnneeded=true
Requires=openproject.service
After=openproject.service

[Service]
# this service is just a placeholder
ExecStart=/bin/sleep infinity

[Install]
WantedBy=openproject.service
