# vendor/pkgr/scaling/systemd/openproject-check-PROCESS_NUM.service
[Unit]
StopWhenUnneeded=true
Requires=openproject-check.service
After=openproject-check.service

[Service]
Environment=APP_PROCESS_INDEX=__PROCESS_INDEX__
ExecStart=/usr/bin/openproject run check
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=%n

[Install]
WantedBy=openproject-check.service
