# vendor/pkgr/scaling/systemd/openproject-worker-PROCESS_NUM.service
[Unit]
StopWhenUnneeded=true
Requires=openproject-worker.service
After=openproject-worker.service

[Service]
Environment=APP_PROCESS_INDEX=__PROCESS_INDEX__
ExecStart=/usr/bin/openproject run worker
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=%n

[Install]
WantedBy=openproject-worker.service
