# vendor/pkgr/scaling/systemd/openproject-console-PROCESS_NUM.service
[Unit]
StopWhenUnneeded=true
Requires=openproject-console.service
After=openproject-console.service

[Service]
Environment=APP_PROCESS_INDEX=__PROCESS_INDEX__
ExecStart=/usr/bin/openproject run console
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=%n

[Install]
WantedBy=openproject-console.service
